# etc/cuirass-remote-worker.service.in
# This is a "service unit file" for the systemd init system to launch
# 'cuirass remote-worker'.  Drop it in /etc/systemd/system or similar
# to have 'cuirass remote-worker' automatically started.

[Unit]
Description=Continuous integration remmote worker for GNU Guix

[Service]
ExecStart=@guix_localstatedir@/guix/profiles/per-user/root/guix-profile/bin/cuirass \
    remote-worker --workers=2 --server=10.0.0.1:5555
Environment= LC_ALL=en_US.UTF-8 \
    'GUIX_LOCPATH=@guix_localstatedir@/guix/profiles/per-user/root/guix-profile/lib/locale'

[Install]
WantedBy=multi-user.target
